# contrib/systemd/ha_cluster/dgraph-alpha.service
[Unit]
Description=dgraph alpha server
Wants=network.target
After=network.target

[Service]
Type=simple
WorkingDirectory=/var/lib/dgraph
Restart=on-failure
ExecStart=/bin/bash -c '/usr/local/bin/dgraph alpha --my={{ myhostname }}:7080 --zero {{ zero-0 }}:5080,{{ zero-1 }}:5080,{{ zero-2 }}:5080 --postings /var/lib/dgraph/p --wal /var/lib/dgraph/w'
StandardOutput=journal
StandardError=journal
User=dgraph
Group=dgraph

[Install]
WantedBy=multi-user.target
